# conf.d/x11vnc.conf
[program:x11vnc]
#command=x11vnc -forever -shared
command=x11vnc -create -forever -shared -loop
autorestart=true